# overlay/system/lib/systemd/system/repowerd.service.d/citrus.conf
[Service]
# Android lights HAL is present, but does not work correctly with repowerd
Environment=REPOWERD_BACKLIGHT_BACKEND=sysfs
Environment=REPOWERD_DISABLE_BOOSTER=true
